# 99-adata-hdd.rules
#
SUBSYSTEMS=="block",KERNEL=="sd[a-z]", RUN+="socket:/org/kernel/udev/adata_hdd"
